# minute_meeting.service
[Unit]
Description=Minute Meeting Transcription Service
After=network.target

[Service]
Type=simple
User=user1
WorkingDirectory=/home/user1
ExecStart=/home/user1/minute_meeting/.venv/bin/python -m uvicorn minute_meeting.app:app --host 0.0.0.0 --port 3318
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target